# etc/supervisor.d/available/results.conf

[program:results]
command=/Users/mw/.virtualenvs/fishfaceenv//bin/celery worker --app=lib.workers.results_tasks -l INFO -Q results -n 'results.%%h' --concurrency=4 -O fair 
directory=/Users/mw/FishFace2
user=mw
numprocs=1
stdout_logfile=/Users/mw/FishFace2/var/log/workers/results.log
stderr_logfile=/Users/mw/FishFace2/var/log/workers/results.err
autostart=true
autorestart=true
stopwaitsecs=600
killasgroup=true
priority=100